# supervisord.conf
[supervisord]
nodaemon=true

#[program:version]
#command=python ./version.py

[program:selenium]
command=java -jar selenium-server-standalone-3.0.1.jar